# dock/gentoobb/images/exim-dovecot/etc/service/exim/run
#!/bin/sh

exec 2>&1
exec /usr/sbin/mysqld --defaults-file=/etc/mysql/my.cnf --user=mysql --log-bin=mysqld-bin
